# Consul/docker/supervisord/consul.conf
[program:nginx]
command=/bin/consul agent -server -bootstrap-expect 3 -data-dir /opt/consul
autostart = true
stdout_events_enabled=true
stderr_events_enabled=true