# Banking/docker/supervisord/external-worker.conf
[program:nginx]
command=php /var/www/src/UI/Worker/External/SynchronizeEvents.php
autostart = true
stdout_events_enabled=true
stderr_events_enabled=true