# Client/docker/supervisord/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart = true
stdout_events_enabled=true
stderr_events_enabled=true