# Client/docker/supervisord/php5-fpm.conf
[program:php5-fpm]
command = /usr/sbin/php5-fpm
autostart = true

[program:php5-fpm-log]
command=tail -f /var/log/php5-fpm.log

